# salt-example/salt-master/supervisor-salt-orig.conf
[program:salt-master]
command=/usr/bin/salt-master
numprocs=1
autostart=false
autorestart=false

[program:sshd]
command=/usr/sbin/sshd -D
numprocs=1
autostart=false
autorestart=false

